# network/env.create.service
[Unit]
Description=Generates /home/kubernetes/env/
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash /home/kubernetes/network/env.create.sh
