# master/kube-master-register.consul.timer
[Unit]
Description=Register kubernetes master every 10sec

[Install]
WantedBy=timers.target

[Timer]
OnBootSec=0min
OnCalendar=*:*:0,10,20,30,40,50